# src/etc/supervisor/conf.d/rsyslogd.conf
[supervisord]
nodaemon=true

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
