# src/etc/supervisor/conf.d/rsyslogd.conf
[supervisord]
nodaemon=true

[program:rsyslog]
command=/bin/bash -c "service rsyslog start"
